# system/ostree-finalize-staged.service.d/rpi-bootc-bootloader.conf
[Service]
ExecStop=
ExecStop=/usr/bin/rpi-bootc-bootloader finalize-staged
